# src/scraper/crontab
10 * * * * scrape-instructors
20 */2 * * * scrape-rate-my-professors
*/6 * * * * scrape-sections
40 * * * * scrape-courses
45 * * * * scrape-transfer-courses
